# docker/scheduler/crontab
# Executa em dias uteis as 10:00 (America/Sao_Paulo).
# O script decide se deve disparar pela regra 5/15/25 + proximo dia util.
0 10 * * 1-5 /scheduler/run-motor.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
